# deploy/systemd/t3-janitor.service
[Unit]
Description=CASM T3 dump-directory janitor (corr1 manages both nodes)
After=network-online.target

[Service]
ExecStart=/home/casm/software/dev/casm_venvs/casm_offline_env/bin/t3-janitor --max-gb 150 --max-age-days 7 --interval-s 3600 --log-file /mnt/nvme5/casm_pipeline/logs/t3_janitor.log
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
