# deploy/systemd/t3-web.service
[Unit]
Description=CASM monitor web UI (port 8050)
After=network-online.target

[Service]
ExecStart=/home/casm/software/dev/casm_venvs/casm_offline_env/bin/t3-web
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
